# solaredge-monitor.service
[Unit]
Description=SolarEdge Off-Grid Monitor
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/solaredge-monitor
ExecStart=/home/pi/solaredge-monitor/venv/bin/python3 main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=200
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=solaredge-monitor

[Install]
WantedBy=multi-user.target
